# luci-app-daqemon/root/etc/hotplug.d/usb/40-daqemon
#!/bin/sh
# This file is autogenerated, do not edit it
/usr/lib/lua/daqemon/daqemon-hotplug.sh /devices/platform/101c1000.ohci/usb2/2-1
